# jan-developer.timer
[Unit]
Description=Run Jan Developer Agent periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1s

[Install]
WantedBy=timers.target
